# crontab/cron.sh
#!/bin/bash


30 0 * * * root find / -name '.DS_Store' -type f -delete
